# EatFitService/supervisor.celerycam.conf
[program:djangoproject.celerycam]
command=/home/eatfit/eatfitservice/eatfitrepo/EatFitService/env/bin/python /home/eatfit/eatfitservice/eatfitrepo/EatFitService/manage.py celerycam --frequency=10.0
user=www-data
numprocs=1
directory=/home/eatfit/eatfitservice/eatfitrepo/EatFitService
stdout_logfile=/home/eatfit/eatfitservice/eatfitrepo/EatFitService/logs/celerycam.log
stderr_logfile=/home/eatfit/eatfitservice/eatfitrepo/EatFitService/logs/celerycam.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 120
priority=998